# files/etc/uci-defaults/enable_UPNP
#!/bin/sh
#Power By Qs315490

# 开启 UPNP
if [ -e /etc/config/upnpd ]
then
        uci set upnpd.config.enabled='1'
        uci commit
fi