# files/etc/uci-defaults/set_bash
#!/bin/sh
#Power By Qs315490

if [ -e /bin/bash ]
then
        # 修改默认 shell
        sed -i 's,/ash,/bash,g' /etc/passwd
fi