# files/etc/uci-defaults/theme_argon
#!/bin/sh
#Power By Qs315490

# 设置主题
if [ -e /www/luci-static/argon ]
then
        uci set luci.main.mediaurlbase='/luci-static/argon'
        uci commit
fi
